# infra/com.pixel-office.reporter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pixel-office.reporter</string>

  <key>ProgramArguments</key>
  <array>
    <string>/path/to/node</string>
    <string>/path/to/pixel-office/pixel-office-reporter.js</string>
    <string>ws://YOUR_SERVER:3300/ws/report</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/pixel-office-reporter.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/pixel-office-reporter.log</string>
</dict>
</plist>
