# infra/pixel-office-display.service
[Unit]
Description=Pixel Office Wall Display Monitor
After=network-online.target pixel-office-server.service
Wants=network-online.target pixel-office-server.service

[Service]
Type=simple
User=root
ExecStart=/home/fedev/homeserver/scripts/pixel-office-display.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Memory limits — MemoryHigh is the soft throttle threshold.
# Too low causes constant kernel throttling (freezes JS event loop).
MemoryMax=2G
MemoryHigh=1800M

[Install]
WantedBy=multi-user.target
