# infra/pixel-office-display.service
[Unit]
Description=Pixel Office Wall Display Monitor
After=network-online.target pixel-office-server.service
Wants=network-online.target pixel-office-server.service

[Service]
Type=simple
User=root
ExecStart=/home/fedev/homeserver/scripts/pixel-office-display.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Memory limits — Chrome kiosk with canvas + WebSocket uses ~1.5GB baseline.
# MemoryHigh is the soft throttle threshold; too low causes D-state freezes
# as the kernel's memory reclaim puts Chrome processes in uninterruptible sleep.
MemoryMax=4G
MemoryHigh=3G

[Install]
WantedBy=multi-user.target
