# infra/pixel-office-server.service
# Adjust User, WorkingDirectory, and ExecStart paths for your system
[Unit]
Description=Pixel Office Server (Node.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<user>
WorkingDirectory=/home/<user>/pixel-office
ExecStart=/usr/bin/node /home/<user>/pixel-office/standalone-server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
