# infra/pixel-office-server.service
[Unit]
Description=Pixel Office Server (Node.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fedev
WorkingDirectory=/home/fedev/pixel-office
ExecStart=/usr/bin/node /home/fedev/pixel-office/standalone-server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
