# infra/pixel-office-display.service
[Unit]
Description=Pixel Office Wall Display Monitor
After=network-online.target pixel-office-server.service
Wants=network-online.target pixel-office-server.service

[Service]
Type=simple
User=root
ExecStart=/home/fedev/homeserver/scripts/pixel-office-display.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Memory limit — hard kill only. Do NOT set MemoryHigh: it triggers kernel
# memory reclaim that puts Chrome in D-state (uninterruptible sleep), freezing
# the render pipeline. Chrome's memory grows ~500MB/min due to GPU compositor
# buffers on AMD Rembrandt; the display script handles this via periodic restart.
MemoryMax=6G

[Install]
WantedBy=multi-user.target
